# packages/github.com/stargrave/goircd/startup/upstart.conf
description "goircd IRC server"

start on (started networking)
stop on (stopped networking)
console log
exec /home/goircd/goircd -hostname irc.example.com
