# misc/vmbetter_configs/generic_vm_overlay/etc/rc.local
#!/bin/sh -e

mount -o remount,rw /
dhclient eth0
exit 0
